# deploy/systemd/anomyze-audit-backup.timer
[Unit]
Description=Daily GPG-encrypted audit-log backup

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=15m
Unit=anomyze-audit-backup.service

[Install]
WantedBy=timers.target
